# systemd/littlebro-motion.service
[Service]
User=www-data
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/littlebro-motion

[Install]
WantedBy=multi-user.target