# systemd/littlebro-archive.service
[Service]
User=www-data
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/littlebro-archive

[Install]
WantedBy=multi-user.target